# webservice.service
[Unit]
Description=Scailx simple python webgui
After=go2rtc.service

[Service]
ExecStart=scailx-webgui

[Install]
WantedBy=multi-user.target
